# image-builder/stage-nasos/02-configure-services/files/systemd/nasos-backend.service
[Unit]
Description=nasOS Backend API Server
After=network-online.target nasos-firstboot.service
Wants=network-online.target

[Service]
Type=simple
User=nasos
Group=nasos
SupplementaryGroups=shadow
WorkingDirectory=/opt/nasos/backend
Environment=NASOS_DEV_MODE=false
Environment=NASOS_DATA_DIR=/opt/nasos/data
Environment=NASOS_DB_PATH=/opt/nasos/data/nasos.db
ExecStart=/opt/nasos/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
# ProtectSystem=full makes /usr and /boot read-only but leaves /etc writable
# (the backend needs sudo for share config: smb.conf, /etc/exports)
# NoNewPrivileges is NOT set — backend needs sudo via share-helper.sh
ProtectSystem=full
ReadWritePaths=/opt/nasos/data /srv/nasos
ReadOnlyPaths=/etc/shadow
PrivateTmp=true

[Install]
WantedBy=multi-user.target
