# image-builder/stage-nasos/02-configure-services/files/systemd/nasos-electron.service
[Unit]
Description=nasOS Desktop (Electron via Cage Wayland compositor)
After=nasos-backend.service nasos-firstboot.service systemd-logind.service
Wants=nasos-backend.service
Conflicts=getty@tty1.service

[Service]
Type=simple
User=nasos
Group=nasos
PAMName=login
SupplementaryGroups=video render input

# Environment for Wayland compositor
Environment=WLR_RENDERER=pixman
# Required for Pi 5 (vc4/v3d DRM driver): wlroots atomic modesetting can fail
# on the Pi 5 display pipeline — disabling it forces legacy modesetting which works.
Environment=WLR_DRM_NO_ATOMIC=1
# Hardware cursors can cause rendering artefacts with the pixman software renderer
Environment=WLR_NO_HARDWARE_CURSORS=1
# NOTE: WLR_LIBINPUT_NO_DEVICES must NOT be set — removing it allows wlroots/cage
# to use libinput, which is required for mouse/touchpad input to work.

# Create the XDG_RUNTIME_DIR for the service user (runs as root via "+" prefix)
# Uses id -u to dynamically resolve the UID even if Pi Imager renamed the user
ExecStartPre=+/bin/bash -c 'UID_VAL=$(id -u nasos 2>/dev/null || echo 1000); mkdir -p /run/user/$UID_VAL; chown $UID_VAL:$UID_VAL /run/user/$UID_VAL; chmod 700 /run/user/$UID_VAL'

# Set XDG_RUNTIME_DIR dynamically via bash wrapper so it resolves the actual UID
# --no-sandbox / --disable-gpu-sandbox: required when Electron runs as a system service
# --disable-renderer-backgrounding etc.: prevent Chromium from throttling/blanking the
#   renderer when it thinks the window is in the background (very common on kiosk setups)
# --disable-features=PaintHolding: prevents blank frames being held during content load
ExecStart=/bin/bash -c 'export XDG_RUNTIME_DIR=/run/user/$(id -u); exec /usr/bin/cage -- /opt/nasos/electron/node_modules/.bin/electron /opt/nasos/electron/main.js --no-sandbox --disable-gpu-sandbox --disable-renderer-backgrounding --disable-background-timer-throttling --disable-backgrounding-occluded-windows --disable-features=PaintHolding'
Restart=on-failure
RestartSec=3

# TTY for display output
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
