# image-builder/stage-nasos/02-configure-services/files/systemd/nasos-firstboot.service
[Unit]
Description=nasOS First Boot Setup
ConditionPathExists=!/opt/nasos/data/.setup-complete
After=local-fs.target network.target
Before=nasos-backend.service nasos-electron.service ssh.service sshd.service

[Service]
Type=oneshot
ExecStart=/opt/nasos/scripts/first-boot.sh
RemainAfterExit=yes
# Safety valve — if something hangs, fail after 5 minutes instead of blocking forever
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
